# vpn-shop/vpn-shop-web.service.example
[Unit]
Description=VPN Shop Web Checkout - SilentConnect
After=network-online.target x-ui.service subjson.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/vpn-shop
EnvironmentFile=/root/vpn-shop/.env
ExecStart=/usr/bin/python3 -m vpn_shop --env-file .env.silentconnect run-web
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
